# cron/gecko-alpha.crontab
# === BEGIN gecko-alpha managed block (do not edit between sentinels) ===
30 3 * * 0 /root/gecko-alpha/scripts/tg_burst_archive.sh >> /var/log/gecko-alpha-tg-burst-archive.log 2>&1
45 3 * * 0 /root/gecko-alpha/scripts/wal_archive.sh >> /var/log/gecko-alpha-wal-archive.log 2>&1
*/5 * * * * /root/gecko-alpha/scripts/source-calls-live-writer.sh >> /var/log/gecko-alpha-source-calls-live-writer.log 2>&1
*/10 * * * * /root/gecko-alpha/scripts/source-calls-lag-watchdog.sh >> /var/log/gecko-alpha-source-calls-lag-watchdog.log 2>&1
*/10 * * * * cd /root/gecko-alpha && .venv/bin/python scripts/check_trade_decision_events.py --db /root/gecko-alpha/scout.db >> /var/log/gecko-alpha-trade-decision-events-watchdog.log 2>&1
20 9 * * * /root/gecko-alpha/scripts/audit_stop_loss_false_negatives.sh --alert /root/gecko-alpha/scout.db >> /var/log/gecko-alpha-stop-loss-fn-audit.log 2>&1
*/15 * * * * /root/gecko-alpha/scripts/acceleration-heartbeat-watchdog.sh >> /var/log/gecko-alpha-acceleration-heartbeat-watchdog.log 2>&1
*/5 * * * * /root/gecko-alpha/scripts/held-position-price-watchdog.sh >> /var/log/gecko-alpha-held-position-price-watchdog.log 2>&1
30 9 * * * /root/gecko-alpha/scripts/revival-verdict-watchdog.sh >> /var/log/gecko-alpha-revival-verdict-watchdog.log 2>&1
15 4 * * 1 /root/gecko-alpha/scripts/ledger-eviction-export.sh >> /var/log/gecko-alpha-ledger-eviction-export.log 2>&1
# alert-channel + digest freshness watchdog (§12a): hourly, monitors tg_alert_log
# ('sent' freshness) + paper_daily_summary (daily digest) in one script. The
# inline ALERT_CHANNEL_WATCHDOG_ENABLED=true is the deploy-without-activate flag
# (safe no-op if unset); ACTIVATION occurs when the operator runs cron/deploy.sh
# to install this managed block, per operator approval.
50 * * * * ALERT_CHANNEL_WATCHDOG_ENABLED=true /root/gecko-alpha/scripts/alert-channel-watchdog.sh >> /var/log/gecko-alpha-alert-channel-watchdog.log 2>&1
# === END gecko-alpha managed block ===
