# cron/gecko-alpha.crontab
# === BEGIN gecko-alpha managed block (do not edit between sentinels) ===
30 3 * * 0 /root/gecko-alpha/scripts/tg_burst_archive.sh >> /var/log/gecko-alpha-tg-burst-archive.log 2>&1
45 3 * * 0 /root/gecko-alpha/scripts/wal_archive.sh >> /var/log/gecko-alpha-wal-archive.log 2>&1
*/5 * * * * /root/gecko-alpha/scripts/source-calls-live-writer.sh >> /var/log/gecko-alpha-source-calls-live-writer.log 2>&1
*/10 * * * * /root/gecko-alpha/scripts/source-calls-lag-watchdog.sh >> /var/log/gecko-alpha-source-calls-lag-watchdog.log 2>&1
*/10 * * * * cd /root/gecko-alpha && .venv/bin/python scripts/check_trade_decision_events.py --db /root/gecko-alpha/scout.db >> /var/log/gecko-alpha-trade-decision-events-watchdog.log 2>&1
20 9 * * * /root/gecko-alpha/scripts/audit_stop_loss_false_negatives.sh --alert /root/gecko-alpha/scout.db >> /var/log/gecko-alpha-stop-loss-fn-audit.log 2>&1
*/15 * * * * /root/gecko-alpha/scripts/acceleration-heartbeat-watchdog.sh >> /var/log/gecko-alpha-acceleration-heartbeat-watchdog.log 2>&1
*/5 * * * * /root/gecko-alpha/scripts/held-position-price-watchdog.sh >> /var/log/gecko-alpha-held-position-price-watchdog.log 2>&1
30 9 * * * /root/gecko-alpha/scripts/revival-verdict-watchdog.sh >> /var/log/gecko-alpha-revival-verdict-watchdog.log 2>&1
15 4 * * 1 /root/gecko-alpha/scripts/ledger-eviction-export.sh >> /var/log/gecko-alpha-ledger-eviction-export.log 2>&1
# alert-channel + digest freshness watchdog (§12a): hourly, monitors tg_alert_log
# ('sent' freshness) + paper_daily_summary (daily digest) in one script. The
# inline ALERT_CHANNEL_WATCHDOG_ENABLED=true is the deploy-without-activate flag
# (safe no-op if unset); ACTIVATION occurs when the operator runs cron/deploy.sh
# to install this managed block, per operator approval. Per-table 24h send
# cooldown bounds a standing breach to one page/table/day (not ~24).
# *** ACTIVATION PREREQUISITE: deploy PR #429 (daily-digest yesterday-fix) and
# confirm >=1 fresh paper_daily_summary row exists BEFORE activating, else the
# first digest pages are for a known-broken-being-fixed writer. ***
50 * * * * ALERT_CHANNEL_WATCHDOG_ENABLED=true /root/gecko-alpha/scripts/alert-channel-watchdog.sh >> /var/log/gecko-alpha-alert-channel-watchdog.log 2>&1
# Liquidity enrichment (b2 Phase 1a-ii). Writer is a NO-OP unless
# LIQUIDITY_ENRICHMENT_ENABLED=true in .env (killswitch). See
# tasks/runbook_enable_liquidity_enrichment_2026_06_23.md for the enable
# sequence — flip the .env flag and verify a manual run BEFORE deploying
# this block, or the watchdog will (correctly) alert writer-heartbeat staleness.
*/5 * * * * cd /root/gecko-alpha && .venv/bin/python scripts/backfill_dexscreener_liquidity.py --db /root/gecko-alpha/scout.db --heartbeat-file /var/lib/gecko-alpha/liquidity-enrichment/heartbeat >> /var/log/gecko-alpha-liquidity-enrichment-writer.log 2>&1
*/10 * * * * cd /root/gecko-alpha && .venv/bin/python scripts/check_liquidity_enrichment_lag.py --db /root/gecko-alpha/scout.db --writer-heartbeat-file /var/lib/gecko-alpha/liquidity-enrichment/heartbeat >> /var/log/gecko-alpha-liquidity-enrichment-watchdog.log 2>&1
# weekly suppression-cost rollup (#432): Monday 09:30, reads the #421 recall lane
# (deployed 2026-07-10). n-gated: prints INSUFFICIENT_DATA until >=10 distinct
# matured tokens (~2026-07-17+); never emits a dollar figure under the gate.
30 9 * * 1 cd /root/gecko-alpha && .venv/bin/python scripts/suppression_cost_rollup.py --db /root/gecko-alpha/scout.db --send >> /var/log/gecko-alpha-suppression-cost-rollup.log 2>&1
# === END gecko-alpha managed block ===
