# systemd/gecko-backup.service
[Unit]
Description=Gecko-Alpha — rotate scout.db backups, keep top-N most-recent

[Service]
Type=oneshot
User=root
Group=root
StateDirectory=gecko-alpha
StateDirectoryMode=0750
Environment=GECKO_BACKUP_DIR=/root/gecko-alpha
Environment=GECKO_BACKUP_KEEP=3
Environment=GECKO_BACKUP_HEARTBEAT_FILE=/var/lib/gecko-alpha/backup-last-ok
ExecStartPre=/usr/bin/test -x /usr/local/bin/gecko-backup-rotate.sh
ExecStart=/usr/local/bin/gecko-backup-rotate.sh
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
