# systemd/minara-emission-persistence-watchdog.service
[Unit]
Description=Gecko-Alpha - Minara emission persistence parity watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/gecko-alpha
Environment=GECKO_APP_DIR=/root/gecko-alpha
Environment=GECKO_DB_PATH=/root/gecko-alpha/scout.db
Environment=GECKO_ENV_FILE=/root/gecko-alpha/.env
Environment=GECKO_PYTHON=/root/gecko-alpha/.venv/bin/python
Environment=MINARA_EMISSION_WATCHDOG_SERVICE=gecko-pipeline
Environment=MINARA_EMISSION_WATCHDOG_SINCE=24 hours ago
Environment=MINARA_EMISSION_WATCHDOG_TOLERANCE=0
ExecStartPre=/usr/bin/test -x /root/gecko-alpha/scripts/minara-emission-persistence-watchdog.sh
ExecStart=/root/gecko-alpha/scripts/minara-emission-persistence-watchdog.sh
# Exit 1 = parity deficit detected and Telegram alert delivered.
# Other nonzero exits mean watchdog diagnostic failure or alert delivery failure.
SuccessExitStatus=1
TimeoutStartSec=90
StandardOutput=journal
StandardError=journal
