# systemd/solana-execution-watchdog.service
[Unit]
Description=gecko-alpha Solana execution freshness watchdog
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/gecko-alpha
EnvironmentFile=/opt/gecko-alpha/.env
ExecStart=/opt/gecko-alpha/.venv/bin/python -m scripts.solana_execution_watchdog
