# systemd/solana-execution-watchdog.timer
[Unit]
Description=Run gecko-alpha Solana execution watchdog hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
