# systemd/solana-sweep.timer
[Unit]
Description=Run gecko-alpha Solana float sweep daily

[Timer]
OnCalendar=*-*-* 00:10:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
